# labvision.service
# This is a reference template. The install.sh script generates the actual
# service file at /etc/systemd/system/labvision.service with the
# correct paths and user filled in automatically.
#
# To install as a service, run:
#   sudo bash install.sh
#
# Manual install (if you prefer):
#   1. Copy this file to /etc/systemd/system/labvision.service
#   2. Replace the placeholder values below
#   3. sudo systemctl daemon-reload
#   4. sudo systemctl enable labvision
#   5. sudo systemctl start labvision

[Unit]
Description=LabVision Monitor Dashboard
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=YOUR_USERNAME
Group=YOUR_USERNAME
WorkingDirectory=/path/to/LabVision
ExecStart=/path/to/LabVision/venv/bin/python app.py
Environment="PYTHONUNBUFFERED=1"

Restart=always
RestartSec=10
TimeoutStopSec=30
KillSignal=SIGTERM

StandardOutput=journal
StandardError=journal
SyslogIdentifier=labvision

NoNewPrivileges=true
 PrivateTmp=true

[Install]
WantedBy=multi-user.target
